# services/wiimote.conf
start on started bluetooth and started udev
exec bash /home/kiosk/.emu/services/wiimote.sh
